# modules/X11/xfce4/dist/.config/autostart/keybase_autostart.desktop
[Desktop Entry]
Name=Keybase
Comment=Keybase Filesystem Service and GUI
Type=Application
Exec=env KEYBASE_AUTOSTART=1 keybase-gui
Hidden=true

