# modules/X11/xfce4/dist/.config/autostart/Toggl.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=Toggl
Comment=Toggl desktop
Exec=/usr/bin/toggldesktop
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false

